# systemd/klipper-gcode-scanner.service
[Unit]
Description=Klipper G-code Scanner
After=network.target

[Service]
Type=simple
WorkingDirectory=@REPO_DIR@
Environment=PYTHONPATH=@REPO_DIR@/src
ExecStart=/usr/bin/python3 -m klipper_gcode_scanner.cli --config /etc/klipper-gcode-scanner/config.toml daemon
Restart=on-failure
User=@USER@
Group=@GROUP@

[Install]
WantedBy=multi-user.target
